# etc/template.service
[Unit]
Description=Tempaast-#{probeName}
After=multi-user.target

[Service]
Type=simple
ExecStart=python #{installation_dir}/main.py "#{probeDir}" "#{probeName}" "#{userId}"
StandardInput=tty-force

[Install]
WantedBy=multi-user.target
